# ops/systemd/netric-nightly-refresh.service
[Unit]
Description=Netric nightly NBA refresh
After=network-online.target
Wants=network-online.target
StartLimitIntervalSec=0

[Service]
Type=oneshot
WorkingDirectory=__NETRIC_BACKEND_DIR__
Environment=PYTHONUNBUFFERED=1
Environment=MONGO_INIT_RETRIES=12
Environment=MONGO_INIT_RETRY_DELAY_SECONDS=10
Environment=MONGO_SERVER_SELECTION_TIMEOUT_MS=15000
Environment=FETCH_WORKER_RETRY_DELAY_SECONDS=900
Environment=NBA_MAX_SEASONS_TO_STORE=1
Environment=NBA_REFRESH_AUDIT_PLAYER_IDS=
ExecStart=/usr/bin/flock -n /tmp/netric-nightly-refresh.lock __NETRIC_PYTHON_BIN__ __NETRIC_BACKEND_DIR__/nightly_refresh.py
Restart=on-failure
RestartSec=5m

[Install]
WantedBy=default.target
