# ops/systemd/netric-nightly-refresh.timer
[Unit]
Description=Run Netric NBA cache refresh checks

[Timer]
OnCalendar=*-*-* 08:00:00
OnBootSec=5m
OnUnitActiveSec=30m
Persistent=true
Unit=netric-nightly-refresh.service

[Install]
WantedBy=timers.target
