# ops/systemd/netric-nightly-refresh.service
[Unit]
Description=Netric nightly NBA refresh
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
WorkingDirectory=/home/dxle/dev/coding-in-class/spring2026/4800CS/netric/backend
ExecStart=/usr/bin/flock -n /tmp/netric-nightly-refresh.lock /home/dxle/dev/coding-in-class/spring2026/4800CS/netric/backend/venv/bin/python /home/dxle/dev/coding-in-class/spring2026/4800CS/netric/backend/nightly_refresh.py

[Install]
WantedBy=default.target
